# pre_config_files/supervisor.conf
[supervisord]
nodaemon=true

[program:solr]
command=/bin/bash start-solr.sh
directory=/usr/local/solr8/pre_config_files
autostart = true
autorestart = false
stopasgroup=true

